# systemd/ai.hermes.seed-registry.timer
[Unit]
Description=Sid Trading Lab — Seed Registry Timer (weekly)

[Timer]
OnCalendar=Sun *-*-* 16:00:00
Persistent=true

[Install]
WantedBy=timers.target
